# docs/systemd-services/freqtrade-bot1.service
[Unit]
Description=FreqTrade Bot 1 (LEA-LSTM Strategy)
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/lea-freqai-system
EnvironmentFile=/home/pi/lea-freqai-system/.env
ExecStart=/home/pi/lea-freqai-system/.venv/bin/freqtrade trade --config /home/pi/lea-freqai-system/user_data/config.json --strategy LeaFreqAIStrategy --freqaimodel PyTorchMLPRegressor --logfile /home/pi/lea-freqai-system/freqtrade.log
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=freqtrade-bot1

[Install]
WantedBy=multi-user.target
